# workshops/systemd/catmorning@.path
[Unit]
Description=Cat on catmorning file changed

[Path]
PathChanged=/home/%i/catmorning.txt
Unit=catmorning@%i.service