# workshops/systemd/catmorning@.service
[Unit]
Description=Cat a file to the service logs

[Service]
Type=oneshot
ExecStart=/bin/bash -c "cat /home/%i/catmorning.txt"